# infra/heimdal/gpu-fan-max.service
[Unit]
Description=Force AMD GPU fan (hwmon2) to PWM max
# Run after amdgpu module is loaded and hwmon is populated.
After=systemd-modules-load.service
ConditionPathExists=/sys/class/hwmon/hwmon2/pwm1
ConditionPathExists=/sys/class/hwmon/hwmon2/pwm1_enable

[Service]
Type=oneshot
RemainAfterExit=yes
# pwm1_enable=1 → manual control; =2 is auto (firmware curve); =0 is off.
# pwm1=204 is 80% of max (target: lower acoustics with no LLM load).
# Raise back to 255 if running a sustained inference workload.
ExecStart=/bin/sh -c 'echo 1 > /sys/class/hwmon/hwmon2/pwm1_enable; echo 204 > /sys/class/hwmon/hwmon2/pwm1'
# On stop: return to firmware auto curve.
ExecStop=/bin/sh -c 'echo 2 > /sys/class/hwmon/hwmon2/pwm1_enable'

[Install]
WantedBy=multi-user.target
