# infra/heimdal/llama-server.container
[Unit]
Description=llama-server (Vulkan) for Phase 0.5 LLM-agent characterization
Wants=network-online.target
After=network-online.target

[Container]
Image=ghcr.io/ggml-org/llama.cpp:server-vulkan
PublishPort=127.0.0.1:8080:8080
Volume=/var/lib/llama-models:/models:Z,ro
EnvironmentFile=/etc/llama-server.env
AddDevice=/dev/dri/renderD128
# Numeric GIDs for host render (105) and video (39) groups.
# The container image lacks named 'render'/'video' groups, so podman
# --group-add must use numeric IDs to avoid container-side group lookup.
GroupAdd=105
GroupAdd=39
Exec=--host 0.0.0.0 --port 8080 --model /models/${LLAMA_MODEL} --parallel ${LLAMA_SLOTS} --n-gpu-layers 999 --ctx-size 8192 --batch-size 512 --no-mmap

[Service]
# Also load env-file at the systemd unit level so ${LLAMA_MODEL}/${LLAMA_SLOTS}
# in Exec= are substituted before podman is invoked. (The [Container]
# EnvironmentFile= only passes vars into the container via --env-file.)
EnvironmentFile=/etc/llama-server.env
Restart=on-failure
RestartSec=5
TimeoutStartSec=120

[Install]
WantedBy=default.target
